# board/raspberrypi2/udev/99-network.rules
KERNEL=="eth*", SUBSYSTEM=="net", ACTION=="add",  RUN+="/sbin/ifup $env{INTERFACE}"
