# system/skeleton/etc/network/if-up.d/rdate
#!/bin/sh

if  [ `rdate -s __NTP_SERVER__` ]; then
	echo "rdate: success"
fi

