# system/skeleton/etc/network/if-up.d/rdate
#!/bin/sh

if  [ `rdate -s ntp.xs4all.nl` ]; then
	echo "rdate: success"
fi

